# templates/program.conf.erb
[program:<%= @name %>]
autostart=<%= @autostart %>
autorestart=<%= @autorestart %>
command=<%= @command %>
user=<%= @user %>
<%- if @numprocs > 1 -%>
numprocs=<%= @numprocs %>
process_name=%(program_name)s_%(process_num)02d
<%- end -%>
